# etc/kingpower.conf
description "superflux upstart script"

respawn
respawn limit 15 5

start on runlevel [2345]
stop on shutdown

setuid mmm1234
exec /bin/sh moerah1381king/launch.sh
